# curator4-1101994313.conf
[program:curator4-1101994313]
command=python3 /home/curator4/curator4-1101994313/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1101994313/test.err.log
stdout_logfile=/home/curator4/curator4-1101994313/test.out.log